# var/etc/supervisord-celery.conf
[program:celery-beat]
command=env/bin/celery beat --app=app -l info -S django --pidfile="var/tmp/celerybeat.pid"
stdout_logfile=var/log/celeryd.beat.log
stderr_logfile=var/log/celeryderr.beat.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=999
pidfile=var/tmp/celerybeat.pid

[program:celery-default]
command=env/bin/celery worker --hostname=default-node@%n --app=app -l info -Q default --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.default.log
stderr_logfile=var/log/celeryderr.default.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
killasgroup=true
priority=990

[program:celery-nmap]
command=env/bin/celery worker --hostname=nmap-node@%n --app=app -l info -Q scan-nmap,monitor-nmap --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.nmap.log
stderr_logfile=var/log/celeryderr.nmap.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
killasgroup=true
priority=990

[program:celery-nessus]
command=env/bin/celery worker --hostname=nessus-node@%n --app=app -l info -Q scan-nessus,monitor-nessus --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.nessus.log
stderr_logfile=var/log/celeryderr.nessus.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-arachni]
command=env/bin/celery worker --hostname=arachni-node@%n --app=app -l info -Q scan-arachni,monitor-arachni --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.arachni.log
stderr_logfile=var/log/celeryderr.arachni.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

;[program:celery-burp]
;command=env/bin/celery worker --hostname=burp-node@%n --app=app -l info -Q scan-burp,monitor-burp --purge --without-mingle --without-gossip -Ofair
;stdout_logfile=var/log/celeryd.burp.log
;stderr_logfile=var/log/celeryderr.burp.log
;autostart=false
;autorestart=false
;startsecs=5
;stopwaitsecs=60
;killasgroup=true
;priority=990

[program:celery-ssllabs]
command=env/bin/celery worker --hostname=ssllabs-node@%n --app=app -l info -Q scan-ssllabs,monitor-ssllabs --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.ssllabs.log
stderr_logfile=var/log/celeryderr.ssllabs.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-owl_dns]
command=env/bin/celery worker --hostname=owl_dns-node@%n --app=app -l info -Q scan-owl_dns,monitor-owl_dns --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.owl_dns.log
stderr_logfile=var/log/celeryderr.owl_dns.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-owl_leaks]
command=env/bin/celery worker --hostname=owl_leaks-node@%n --app=app -l info -Q scan-owl_leaks,monitor-owl_leaks --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.owl_leaks.log
stderr_logfile=var/log/celeryderr.owl_leaks.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-owl_code]
command=env/bin/celery worker --hostname=owl_code-node@%n --app=app -l info -Q scan-owl_code,monitor-owl_code --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.owl_code.log
stderr_logfile=var/log/celeryderr.owl_code.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-virustotal]
command=env/bin/celery worker --hostname=virustotal-node@%n --app=app -l info -Q scan-virustotal,monitor-virustotal --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.virustotal.log
stderr_logfile=var/log/celeryderr.virustotal.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-urlvoid]
command=env/bin/celery worker --hostname=urlvoid-node@%n --app=app -l info -Q scan-urlvoid,monitor-urlvoid --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.urlvoid.log
stderr_logfile=var/log/celeryderr.urlvoid.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-censys]
command=env/bin/celery worker --hostname=censys-node@%n --app=app -l info -Q scan-censys,monitor-censys --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.censys.log
stderr_logfile=var/log/celeryderr.censys.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[program:celery-cortex]
command=env/bin/celery worker --hostname=cortex-node@%n --app=app -l info -Q scan-cortex,monitor-cortex --purge --without-mingle --without-gossip -Ofair
stdout_logfile=var/log/celeryd.cortex.log
stderr_logfile=var/log/celeryderr.cortex.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=990

[group:celery-workers]
programs=celery-nmap,celery-nessus,celery-ssllabs,celery-arachni,celery-owl_dns,celery-owl_leaks,celery-owl_code,celery-virustotal,celery-urlvoid,celery-censys,celery-cortex
priority=990

[program:celery-flow]
command=env/bin/celery --app=app flower --port=5555
stdout_logfile=var/log/celeryd.flower.log
stderr_logfile=var/log/celeryderr.flower.log
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
killasgroup=true
priority=998